# dot_config/systemd/user/focus-helper.service
[Unit]
Description=Focus Helper - Screen lock after usage threshold
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 %h/.local/share/focus-helper/focus_helper.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
