# dot_config/systemd/user/less-is-more.timer
[Unit]
Description=Run less-is-more command every hour at **:55

[Timer]
OnCalendar=*-*-* *:55:00

[Install]
WantedBy=default.target
