# dot_config/systemd/user/less-is-more.service
[Unit]
Description=Run less-is-more background command

[Service]
ExecStart=/home/diego/prog/all/less-is-more/mk-background.py
